# tmp/droid-config.files
%defattr(-,root,root,-)
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/systemd/system/droid-late-start.service
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-load-firmware.sh
/etc/sap.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/lowi.conf
/etc/izat.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/obexd/noplugins/syncevolution
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/ftp
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/pbap.conf
/etc/obexd/opp.conf
/etc/obexd/root
/etc/media_codecs_google_telephony.xml
/etc/media_profiles.xml
/etc/audio_platform_info.xml
/etc/modprobe.d/blacklist-evbug.conf
/etc/media_codecs_google_audio.xml
/etc/media_codecs_google_video.xml
/etc/xtwifi.conf
/etc/media_codecs_performance.xml
/etc/acdbdata
/etc/mixer_paths.xml
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/binder.conf
/etc/gps.conf
/etc/media_codecs.xml
/etc/media_codecs_ffmpeg.xml
/etc/xtra_root_cert.pem
/etc/selinux/config
/etc/sensorfw/primaryuse.conf
/etc/sec_config
/lib/udev/rules.d/998-droid-system.rules
/lib/udev/rules.d/997-droid-drm.rules
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/ofono/noplugin.conf
/var/lib/environment/nemo/60-multimedia.conf
/usr/libexec/droid-hybris/system/etc/init/drmserver.rc
/usr/libexec/droid-hybris/system/etc/init/rild.rc
/usr/libexec/droid-hybris/system/etc/init/logd.rc
/usr/libexec/droid-hybris/system/etc/init/installd.rc
/usr/libexec/droid-hybris/system/etc/init/init-debug.rc
/usr/libexec/droid-hybris/system/etc/init/debuggerd.rc
/usr/libexec/droid-hybris/system/etc/init/atrace.rc
/usr/libexec/droid-hybris/system/etc/init/dumpstate.rc
/usr/libexec/droid-hybris/system/etc/init/racoon.rc
/usr/libexec/droid-hybris/system/etc/init/mediacodec.rc
/usr/libexec/droid-hybris/system/etc/init/uncrypt.rc
/usr/libexec/droid-hybris/system/etc/init/superuser.rc
/usr/libexec/droid-hybris/system/etc/init/mediaextractor.rc
/usr/libexec/droid-hybris/system/etc/init/logcatd.rc
/usr/libexec/droid-hybris/system/etc/init/bootstat.rc
/usr/libexec/droid-hybris/system/etc/init/mtpd.rc
/usr/libexec/droid-hybris/system/etc/init/mediadrmserver.rc
/usr/libexec/droid-hybris/system/etc/init/vdc.rc
/usr/libexec/droid-hybris/system/etc/init/perfprofd.rc
/usr/libexec/droid-hybris/system/etc/init/gatekeeperd.rc
/usr/libexec/droid-hybris/system/etc/init/vold.rc
/var/lib/environment/compositor/droid-hal-device.conf
/usr/share/ssu/features.d/adaptation-community.ini
